# debian/certbot.service
[Unit]
Description=Certbot

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot -q renew
PrivateTmp=true

